# deploy/cfb-rankings.service
# Systemd service file for College Football Rankings
# Place this in /etc/systemd/system/cfb-rankings.service

[Unit]
Description=College Football Rankings - Gunicorn/FastAPI
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cfb-rankings

# Set environment variables
Environment="PATH=/var/www/cfb-rankings/venv/bin"
Environment="CFBD_API_KEY=EAi7pR9hLKENNyAiU+VdIKalPgYeNHob29XIY18M6QKoYJhrtne8ZVX4u6Oy3d0A"
Environment="DATABASE_URL=sqlite:///./cfb_rankings.db"

# Or use EnvironmentFile if .env exists
# EnvironmentFile=-/var/www/cfb-rankings/.env

ExecStart=/var/www/cfb-rankings/venv/bin/gunicorn -c gunicorn_config.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
