# deploy/cfb-rankings.service
# Systemd service file for College Football Rankings
# Place this in /etc/systemd/system/cfb-rankings.service

[Unit]
Description=College Football Rankings - Gunicorn/FastAPI
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cfb-rankings

# Set environment variables
Environment="PATH=/var/www/cfb-rankings/venv/bin"

# Load secrets (CFBD_API_KEY, DATABASE_URL, CFBD_MONTHLY_LIMIT) from .env
EnvironmentFile=-/var/www/cfb-rankings/.env

ExecStart=/var/www/cfb-rankings/venv/bin/gunicorn -c gunicorn_config.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
