# deploy/cfb-weekly-update.timer
[Unit]
Description=CFB Rankings Weekly Update Timer
Documentation=https://github.com/your-repo/cfb-rankings
Requires=cfb-weekly-update.service

[Timer]
# Trigger every Sunday at 5:00 PM Eastern Time
#
# IMPORTANT: Systemd timers use UTC time. Convert from ET to UTC:
#   - During EDT (Mar-Nov): 5:00 PM ET = 9:00 PM UTC Sunday (21:00)
#   - During EST (Nov-Feb): 5:00 PM ET = 10:00 PM UTC Sunday (22:00)
#
# Using two timer entries to handle DST automatically:
OnCalendar=Sun *-03..10-* 21:00  # EDT months (March-October): 5pm ET = 9pm UTC
OnCalendar=Sun *-11..02-* 22:00  # EST months (November-February): 5pm ET = 10pm UTC

# Persistent=true ensures missed runs (e.g., system was off) will execute on next boot
Persistent=true

# Wait a bit after boot before running (gives network time to stabilize)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
