# deploy/cfb-weekly-update.timer
[Unit]
Description=CFB Rankings Weekly Update Timer
Documentation=https://github.com/your-repo/cfb-rankings
Requires=cfb-weekly-update.service

[Timer]
# Trigger every Sunday at 5:00 PM Eastern Time
#
# IMPORTANT: Systemd timers use UTC time. Convert from ET to UTC:
#   - During EST (Nov-Mar): 5:00 PM ET = 10:00 PM UTC Sunday (22:00)
#   - During EDT (Mar-Nov): 5:00 PM ET = 9:00 PM UTC Sunday (21:00)
#
# Using 22:00 UTC covers EST period (most of CFB season: Sep-Jan)
# During EDT period (late Aug-early Sep), timer runs at 6:00 PM ET instead of 5:00 PM
# This is acceptable as most of the season is played during EST
OnCalendar=Sun 22:00

# Persistent=true ensures missed runs (e.g., system was off) will execute on next boot
Persistent=true

# Wait a bit after boot before running (gives network time to stabilize)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
