# deploy/cfb-rankings.service
# Systemd service file for College Football Rankings
# Place this in /etc/systemd/system/cfb-rankings.service

[Unit]
Description=College Football Rankings - Gunicorn/FastAPI
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=cfb-rankings
WorkingDirectory=/var/www/cfb-rankings
Environment="PATH=/var/www/cfb-rankings/venv/bin"
EnvironmentFile=/var/www/cfb-rankings/.env
ExecStart=/var/www/cfb-rankings/venv/bin/gunicorn -c gunicorn_config.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
